# deploy/secretvoip-sms.service
[Unit]
Description=SecretVoIP SMS backend (Node/Express)
After=network.target postgresql.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/secretvoip-sms/backend
EnvironmentFile=/var/www/secretvoip-sms/.env
ExecStart=/usr/bin/node src/server.js
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/log/secretvoip-sms.log
StandardError=append:/var/log/secretvoip-sms.err.log

[Install]
WantedBy=multi-user.target
